# stu500-tablet.service

[Unit]
Description=STU500 userlevel driver

[Service]
ExecStart=/usr/bin/python3 /usr/local/bin/stu500-tablet.py

